# n8n_monitor.service
[Unit]
Description=N8N Health Monitor with Auto-Restart
After=network.target n8n.service
Wants=n8n.service

[Service]
Type=simple
User=root
# ВАЖНО: Замените /path/to/project на реальный путь к проекту
WorkingDirectory=/path/to/project
EnvironmentFile=/path/to/project/.env
ExecStart=/path/to/project/.venv/bin/python /path/to/project/n8n_monitor.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
